# ecos-tg-bot-ovpn.service
[Unit]
Description=OpenVPN Telegram Management Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/u01/ovpn/ecos-tg-bot-ovpn
EnvironmentFile=/u01/ovpn/ecos-tg-bot-ovpn/.env
ExecStart=/u01/ovpn/ecos-tg-bot-ovpn/venv/bin/python -m bot.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ecos-tg-bot-ovpn

[Install]
WantedBy=multi-user.target
